# skills/x-scheduler/com.miyabi.x-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miyabi.x-scheduler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/.config/last30days/.env && last30days "AI agents autonomous 2026" --quick --json > /tmp/x-research.json && openclaw agent message sns-creator "以下のトレンドから投稿5本生成: " && npx agent-skill-bus record-run --agent sns-creator --skill x-scheduler --task "daily-auto" --result success --score 0.9</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/x-scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/x-scheduler.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
